# etc/skel/systemd/system/disable-pcspkr.service
[Unit]
Description=Alway disable pcspkr

[Service]
ExecStart=rmmod pcspkr
Type=simple

[Install]
WantedBy=default.target
